# humext/openrasp/humsys/supervisord.d/elastic.ini
[program:elastic]
directory=%(ENV_PID_PATH)s
user=root
environment=HOME="/root", USER="root"
priority=99
autostart=false
autorestart=false
startsecs=60
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=100
stdout_capture_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=100
stderr_capture_maxbytes=10MB
stdout_logfile=%(ENV_LOG_PATH)s/stdout/elastic.log
stderr_logfile=%(ENV_LOG_PATH)s/stderr/elastic.log
command=/usr/local/sbin/elastic