# humext/openrasp/humsys/supervisord.d/onedev.ini
[program:onedev]
directory=%(ENV_PID_PATH)s
user=root
environment=HOME="/root", USER="root"
priority=999
autostart=false
autorestart=false
startsecs=25
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=100
stdout_capture_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=100
stderr_capture_maxbytes=10MB
stdout_logfile=%(ENV_LOG_PATH)s/stdout/onedev.log
stderr_logfile=%(ENV_LOG_PATH)s/stderr/onedev.log
command=/usr/local/sbin/onedev