# humsys/supervisord.d/hub.ini
[program:hub]
directory=%(ENV_PID_PATH)s
user=root
environment=HOME="/root", USER="root"
priority=9
autostart=false
autorestart=false
startsecs=15
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=100
stdout_capture_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=100
stderr_capture_maxbytes=10MB
stdout_logfile=%(ENV_LOG_PATH)s/stdout/hub.log
stderr_logfile=%(ENV_LOG_PATH)s/stderr/hub.log
command=/usr/local/bin/jupyterhub --config=%(ENV_HUB_HOME)s/config.py